# kokanybot.service
[Unit]
Description=Kokanybot program

[Service]
ExecStart=/usr/bin/kokanybot

[Install]
WantedBy=multi-user.target
